# systemd/greenLightAdapter.target
[Unit]
Description=greenLightAdapter
[Install]
WantedBy=multi-user.target
